# deploy/echo-research.service
# systemd unit（U-3 / E14）。放到 /etc/systemd/system/echo-research.service，
# 然后 systemctl daemon-reload && systemctl enable --now echo-research。
# .env（API key）放在 WorkingDirectory 下，由应用自己的 loadEnvFile 读取——
# 不写进本文件（红线 20：密钥只存服务器 .env）。

[Unit]
Description=Echo Research (Node + SQLite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=echo
WorkingDirectory=/opt/echo-research
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3
# 生产环境标记：信任 Caddy 的 X-Forwarded-For（限速按真实客户端 IP），cookie 带 Secure
Environment=ECHO_TRUST_PROXY=1
Environment=ECHO_BETA_MODE=1
Environment=NODE_ENV=production
# 资源护栏（2C4G VPS 够用；SQLite 单进程，不需要更多）
MemoryMax=1500M
# 最小权限
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
UMask=0077
LimitNOFILE=4096
ReadWritePaths=/opt/echo-research

[Install]
WantedBy=multi-user.target
